# crontab
* * * * * root cd /app && python sort.py -s /source -d /destination >> /var/log/cron.log 2>&1