# crontab
*/15 * * * * root cd /app && python script.py -s source -d destination >> /var/log/cron.log 2>&1